# docker/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:sshd]
command=/usr/sbin/sshd -D

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman -n

[program:xrdp]
command=/usr/sbin/xrdp -n
